# gpspanel.service
[Unit]
Description=GPS Panel
After=multi-user.target gpsd.socket

[Service]
Type=idle
User=nobody
ExecStart=/usr/bin/python3 /var/www/gpspanel/gpspanel.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
